# ops/phase2/systemd/augur-vllm.service
[Unit]
Description=augur optional vLLM OpenAI server (:8000)
Documentation=file:///home/giga/augur/ops/phase2/VLLM-GB10.md
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/giga/augur
# venv 須含 ninja（Qwen3-4B 編譯路徑需要）；PATH 置前
Environment=PATH=/home/giga/augur-venvs/vllm/bin:/usr/bin
# P1 選定：Qwen/Qwen3-4B／qwen3-4b（見 VLLM-OPT-20260722.md）
ExecStart=/home/giga/augur-venvs/vllm/bin/vllm serve Qwen/Qwen3-4B \
  --host 127.0.0.1 \
  --port 8000 \
  --served-model-name qwen3-4b \
  --max-model-len 8192 \
  --gpu-memory-utilization 0.35
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
# 預設不要 enable；開機自啟須另拍板
